# scripts/odyssey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>in.sixeleven.odyssey</string>

    <!-- Points at an installed copy outside ~/Documents so launchd isn't blocked by macOS TCC.
         Source of truth is scripts/start.sh; re-copy to ~/.local/bin/odyssey-start on change. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/rishitv/.local/bin/odyssey-start</string>
    </array>

    <!-- Run at every login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Don't restart the launcher itself; compose handles container restarts -->
    <key>KeepAlive</key>
    <false/>

    <!-- Watchdog: re-run every 5 min. Script is idempotent and silent when healthy,
         so this just restarts the stack if podman machine or the container dies. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/rishitv/Library/Logs/odyssey.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rishitv/Library/Logs/odyssey-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/rishitv</string>
    </dict>
</dict>
</plist>
